# ops/systemd/cascade-risex-arbi.service
[Unit]
Description=Cascade RISEx arbitrage bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/cascade-risex-arbi
Environment=NODE_ENV=production
EnvironmentFile=/home/ec2-user/cascade-risex-arbi/.env
ExecStart=/usr/bin/node src/main.js
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
